# systemd/system/jenkinshield.socket
[Socket]
ListenStream=2323

[Install]
WantedBy=sockets.target
